# supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisord.log
logfile_maxbytes=10MB
logfile_backups=2
loglevel=info
pidfile=/var/run/supervisord.pid

; ── Main IPTV Manager Server (port 10000) ──────────────────────────────────────
[program:main-server]
command=node /app/server.cjs
directory=/app
autostart=true
autorestart=true
startretries=5
startsecs=5
stopwaitsecs=10
stdout_logfile=/var/log/main-server.log
stdout_logfile_maxbytes=5MB
stdout_logfile_backups=2
stderr_logfile=/var/log/main-server-err.log
stderr_logfile_maxbytes=5MB
stderr_logfile_backups=2
environment=
    PORT="10000",
    NODE_ENV="production",
    DB_FILE="/data/db/db.json",
    STREAM_SERVER_URL="http://localhost:10001"

; ── DRM Streaming Server (port 10001) ──────────────────────────────────────────
[program:streaming-server]
command=node /app/streaming-server/server.js
directory=/app/streaming-server
autostart=true
autorestart=true
startretries=5
startsecs=10
stopwaitsecs=15
stdout_logfile=/var/log/streaming-server.log
stdout_logfile_maxbytes=5MB
stdout_logfile_backups=2
stderr_logfile=/var/log/streaming-server-err.log
stderr_logfile_maxbytes=5MB
stderr_logfile_backups=2
environment=
    STREAM_PORT="10001",
    NODE_ENV="production",
    MAIN_SERVER_URL="http://localhost:10000",
    OUTPUT_DIR="/data/hls-output",
    DB_FILE="/data/db/db.json",
    VIDEO_CODEC="copy",
    AUDIO_CODEC="copy",
    HLS_SEGMENT_TIME="4",
    HLS_LIST_SIZE="6",
    HW_ACCEL="none"
